# linux-ubuntu/fonte/plot-board/opt/openboard/library/applications/plot-board.wgt/plot-board.service
[Unit]
Description = Iniciando o servidor node.js para atender aos requisitos do plot-board.

[Service]
ExecStart=/usr/local/bin/plot-board.sh start

[Install]
WantedBy=default.target